# pi-guy-dashboard.service
[Unit]
Description=Pi-Guy Dashboard
After=network.target

[Service]
Type=simple
User=%u
Environment=PIGUY_HOME=%h/piguy
WorkingDirectory=%h/piguy
ExecStart=/usr/bin/env bash -lc 'exec "${PIGUY_HOME}/start.sh"'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
